# systemd/telegram-codex-bot.service
[Unit]
Description=Telegram Codex Server Bot

[Service]
Type=simple
WorkingDirectory=%h/telegram-codex-bot/telegram-codex-server-bot
EnvironmentFile=%h/telegram-codex-bot/telegram-codex-server-bot/.env
ExecStart=%h/telegram-codex-bot/telegram-codex-server-bot/.venv/bin/python %h/telegram-codex-bot/telegram-codex-server-bot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
